# docker_instances/070_wochenmarkt/crontab
# Run every 10 minutes to keep vendor list up to date
0 * */10 * * . /app/env.sh ; cd /app ; python3 /app/script.py >> /proc/1/fd/1 2>&1
