# docker_instances/070_wochenmarkt/crontab
# Run every 2 hours to keep vendor list up to date
0 */2 * * * . /app/env.sh ; cd /app ; python3 /app/script.py >> /proc/1/fd/1 2>&1
